# conf/dragino_lps8n/files/etc/uci-defaults/90_config_lan
uci -q batch << EOI
set network.lan.proto='dhcp'
delete network.lan.ipaddr
delete network.lan.netmask
delete network.lan.ip6assign
EOI
