# conf/base_raspberrypi_bcm27xx_bcm2709/files/etc/uci-defaults/99_set_chirpstack_mqtt_forwarder_global_config
uci get chirpstack-mqtt-forwarder.@global[0]
if [ $? -eq 1 ]; then
  uci add chirpstack-mqtt-forwarder global
fi
