# conf/rak_rak7391/files/etc/uci-defaults/02_rtc
. /lib/functions/uci-defaults.sh
. /lib/functions.sh
. /lib/functions/system.sh

sed "s|exit 0|echo pcf8563 0x51 > /sys/bus/i2c/devices/i2c-1/new_device\nexit 0|" -i /etc/rc.local 

exit 0