# conf/rak_rak7391/files/etc/uci-defaults/99_set_hostname
if [ "$(uci -q get system.@system[0].hostname)" != "OpenWrt" ]; then
    exit 0
fi

GWID_END=""

ip link show eth0
if [ $? -eq 0 ]; then
    GWID_END=$(ip link show eth0 | awk '/ether/ {print $2}' | awk -F\: '{print $4$5$6}')
else
    GWID_END=$(ip link show wlan0 | awk '/ether/ {print $2}' | awk -F\: '{print $4$5$6}')
fi

uci set system.@system[0].hostname="chirpstack-${GWID_END}"
uci commit system
echo $(uci get system.@system[0].hostname) > /proc/sys/kernel/hostname
