# conf/rak_rak7391/files/etc/uci-defaults/99_set_postgresql_data_path
PG_VERSION=$(postgres -V | awk '{print $NF}')

mkdir -p /srv/postgresql/${PG_VERSION}
chown postgres:postgres /srv/postgresql/${PG_VERSION}

uci set postgresql.config.PGDATA=/srv/postgresql/${PG_VERSION}/data
uci commit
