# conf/base_raspberrypi_bcm27xx_bcm2708/files/etc/uci-defaults/99_set_chirpstack_udp_forwarder_global_config
uci get chirpstack-udp-forwarder.@global[0]
if [ $? -eq 1 ]; then
  uci add chirpstack-udp-forwarder global
fi
