# conf/rak_rak7391/files/etc/rc.local
echo pcf8563 0x51 > /sys/bus/i2c/devices/i2c-1/new_device
exit 0
