# conf/rak_rak7391/files/etc/rc.local
# Import device repository
while :
do
	OUT=$(cmp /opt/lorawan-devices/SRCREV /opt/lorawan-devices/SRCREV.lastimport)
	RET=$?
	if [ $RET -eq 0 ]; then
		break
	fi

	chirpstack -c /etc/chirpstack import-legacy-lorawan-devices-repository -d /opt/lorawan-devices
	RET=$?
	if [ $RET -eq 0 ]; then
		cp /opt/lorawan-devices/SRCREV /opt/lorawan-devices/SRCREV.lastimport
		break
	fi

	sleep 1
done
