# conf/rak_rak7289v2/files/etc/rc.local
# Make gateway-id show first slot Gateway ID
ln -s /tmp/concentratord_slot1_command /tmp/concentratord_command
ln -s /tmp/concentratord_slot1_event /tmp/concentratord_event
