# conf/rak_rak7267/files/etc/hotplug.d/net/10-set-wifi-mac
[ "$ACTION" == "add" ] || exit 0
[ "$DEVTYPE" == "wlan" ] || exit 0

. /lib/functions.sh
. /lib/functions/system.sh

uci set wireless.@wifi-iface[0].macaddr=$(macaddr_add "$(mtd_get_mac_binary factory 0x4)" 1)
uci commit wireless
wifi reload

