# conf/full_raspberrypi_bcm27xx_bcm2709/files/etc/uci-defaults/01_update_rc_local_20241031
# update rc.local file as sysupgrades will preserve the old file
cp /rom/etc/rc.local /etc/rc.local
