# conf/full_raspberrypi_bcm27xx_bcm2709/files/etc/uci-defaults/99_set_postgresql_data_path
mkdir -p /srv/postgresql
chown postgres:postgres /srv/postgresql

uci set postgresql.config.PGDATA=/srv/postgresql/data
uci commit
