# conf/rak_rak7289v2/files/etc/uci-defaults/90_gpsd_config
uci -q batch << EOI
set gpsd.core.enabled='1'
set gpsd.core.device='/dev/ttyS1'
commit gpsd
EOI
