# conf/rak_rak7267/files/etc/uci-defaults/90_migrate_20241008_concentratord_gpsd
uci -q batch << EOI
set chirpstack-concentratord.@sx1302[0].gnss_dev_path='gpsd://localhost:2947'
commit concentratord
EOI
